# scripts/self-healing/com.siriusos.watchdog.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.siriusos.watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{HOME}/.siriusos/{INSTANCE}/scripts/watchdog.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{HOME}/.siriusos/{INSTANCE}/logs/watchdog.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{HOME}/.siriusos/{INSTANCE}/logs/watchdog.stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>CTX_INSTANCE_ID</key>
        <string>{INSTANCE}</string>
    </dict>
</dict>
</plist>
